# assignment4/tarecho.socket
# tarecho.socket

[unit]
Description = Tar server

[Socket]
ListenStream = /var/run/tarecho
Accept = yes

[Install]
WantedBy = sockets.target