# assignment4/tarecho@.service
[Unit]
Description = Assignment #4

[Service]
Type=simple
ExecStart=/usr/local/bin/compile_troff.py